# docker/files/supervisor/celerybeat.conf
[program:celerybeat]
command=celery beat --app pendulum.app --loglevel=WARN
directory=/code

autostart=true
autorestart=true

stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
stderr_logfile=NONE

# Fixed 'missing' log messages, by telling Python to stop buffering the output stream
# (see https://github.com/Supervisor/supervisor/issues/13 for details)
environment=PYTHONUNBUFFERED=1

# How long the program needs to stay running after a startup to
# consider the start successful.
startsecs=10

# How long to wait for the program to exit cleanly, before killing it
stopwaitsecs=15

# Start after RabbitMQ, Redis and Celery
priority=900
