# docker/files/supervisor/rabbitmq.conf
[program:rabbitmq]
command=rabbitmq-supervisor

autostart=true
autorestart=true

stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
stderr_logfile=NONE

# How long the program needs to stay running after a startup to
# consider the start successful.
startsecs=10

# How long to wait for the program to exit cleanly, before killing it
stopwaitsecs=15

# Start RabbitMQ before Celery
priority=100
